# FV3GFSwfm/rt_ccpp_c384/rt_ccpp_c384.crontab

#################### rt_ccpp_c384 ####################
MAILTO=""
*/5 * * * * /apps/rocoto/1.3.3/bin/rocotorun -d /scratch1/BMC/gsd-fv3/rtruns/CCPP-Chem/FV3GFSwfm/rt_ccpp_c384/rt_ccpp_c384.db -w /scratch1/BMC/gsd-fv3/rtruns/CCPP-Chem/FV3GFSwfm/rt_ccpp_c384/rt_ccpp_c384.xml
#################################################################

